# xray_server.service@openwrt
#!/bin/sh /etc/rc.common
START=95

CONF_FOLDER="/etc/xray_server"
RUN_FOLDER="/var/run/xray_server"
BIN_LOG_FOLDER="/var/log/xray_server"

do_start() {
	[ -d $CONF_FOLDER/ ] || return 2
	json_files=$(ls $CONF_FOLDER/ | grep ".json$" )
	[ -d $RUN_FOLDER/ ] && do_stop
	mkdir -p $RUN_FOLDER/
	mkdir -p $BIN_LOG_FOLDER/
	for json_file in $json_files; do
		pre_name=$(echo $json_file | sed "s/\.json//g")
		nohup xray run -c $CONF_FOLDER/$json_file > /dev/null 2>&1 &
		echo $! > $RUN_FOLDER/$pre_name.pid
	done
}

do_stop() {
	[ -d $RUN_FOLDER/ ] || return 2
	pid_files=$(ls $RUN_FOLDER/ | grep ".pid" )
	[ -z "$pid_files" ] && return 2
	for pid_file in $pid_files; do
		pid=$(cat $RUN_FOLDER/$pid_file)
		kill $pid
	done
	rm -rf $RUN_FOLDER/
	rm -rf $BIN_LOG_FOLDER/
}

start() {
	do_start
}

stop() {
	do_stop
}

restart() {
	do_stop
	do_start
}
